# config/dotconfig/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland graphical session
Wants=graphical-session.target
After=graphical-session-pre.target
